# dist/piexpchair_videoplayer.service
[Unit]
Description=PiExperience Chair videoplayer module

[Service]
ExecStart=/home/pi/.local/bin/python3 /home/pi/PiExperienceChair/piexpchair/videoplayer.py
WorkingDirectory=/home/pi/PiExperienceChair/
Restart=always
RestartSec=5
KillSignal=SIGTERM
Type=simple
StandardOutput=journal
StandardError=journal
SyslogIdentifier=piexpchair_videoplayer
NotifyAccess=all
User=pi
Group=pi
Nice=1

[Install]
WantedBy=multi-user.target
