# dist/piexpchair_i2c.service
[Unit]
Description=PiExperience Chair i2c module

[Service]
ExecStart=/home/pi/.local/bin/python3 /home/pi/PiExperienceChair/src/i2c.py
WorkingDirectory=/home/pi/PiExperienceChair/
Restart=always
RestartSec=5
KillSignal=SIGTERM
Type=simple
StandardOutput=journal
StandardError=journal
SyslogIdentifier=piexpchair_i2c
NotifyAccess=all
User=pi
Group=pi
Nice=1

[Install]
WantedBy=multi-user.target
